# systemd/ssrf-health.timer
[Unit]
Description=Run SSRF Console Health Scoring Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
