# systemd/ssrf-health.service
[Unit]
Description=SSRF Console Daily Health Scoring
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/kali/ssrf-command-console
ExecStart=/usr/bin/make health
User=kali
Group=kali
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=false
PrivateTmp=true
